# arch-auto-update.service
[Unit]
Description=Arch Linux Automatic Update Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/arch-auto-update
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target